# systemd/client/qualvosec.service
[Unit]
Description=QualvoSec Service

[Service]
ExecStart=/usr/sbin/qualvosec -c /etc/qualvosec/qualvosec.conf
User=qualvosec
