# systemd/server/qualvosec_server.service
[Unit]
Description=QualvoSec Server Service

[Service]
ExecStart=/usr/sbin/qualvosec_server -c /etc/qualvosec/qualvosec_server.conf
User=qualvosec
